# sv/stalonetray/run
#!/bin/sh
while [ ! -e /tmp/.X99-lock ]; do sleep 0.2; done
exec stalonetray --geometry 5x1+0+0 --icon-size 24 --background "#2a2a2a" --grow-gravity E
