# sv/xvfb/run
#!/bin/sh
rm -f /tmp/.X99-lock
exec Xvfb :99 -screen 0 1024x768x24
